# config/live-build/includes.chroot/usr/lib/systemd/system/nmos-live-user-password.service
[Unit]
Description=Configure the NM-OS live user password before GDM
After=live-config.service systemd-tmpfiles-setup.service
Before=display-manager.service
ConditionPathExists=/usr/local/lib/nmos/ensure_live_user_password.py

[Service]
Type=oneshot
ExecStart=/usr/local/lib/nmos/ensure_live_user_password.py
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
