# config/live-build/includes.chroot/etc/systemd/system/display-manager.service.d/nmos-live-user-password.conf
[Unit]
Wants=nmos-live-user-password.service
After=nmos-live-user-password.service
