# config/system-overlay/usr/lib/systemd/system/nmos-network-ready.target
[Unit]
Description=NM-OS network ready
Documentation=https://github.com/Krypera/NM-OS

