# config/system-overlay/usr/lib/systemd/system/nmos-network-bootstrap.service
[Unit]
Description=NM-OS network bootstrap gate
After=nmos-settings-bootstrap.service
Wants=nmos-settings-bootstrap.service

[Service]
Type=simple
ExecStart=/usr/local/lib/nmos/network_bootstrap.py
Restart=on-failure
RestartSec=5
TimeoutStartSec=320
StandardOutput=journal+console
StandardError=journal+console
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/run/nmos
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK

[Install]
WantedBy=multi-user.target
