# config/system-overlay/usr/lib/systemd/system/nmos-settings.service
[Unit]
Description=NM-OS settings service
After=local-fs.target dbus.service systemd-tmpfiles-setup.service
Wants=dbus.service
Before=display-manager.service

[Service]
Type=simple
ExecStart=/usr/local/bin/nmos-settings-service
Restart=on-failure
RestartSec=3
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=@NMOS_RUNTIME_DIR@ @NMOS_STATE_DIR@
CapabilityBoundingSet=
RestrictAddressFamilies=AF_UNIX

[Install]
WantedBy=multi-user.target
