# config/system-overlay/usr/lib/systemd/system/nmos-logging-policy.service
[Unit]
Description=NM-OS logging policy enforcement
After=nmos-settings-bootstrap.service systemd-journald.service
Wants=nmos-settings-bootstrap.service
Before=nmos-network-bootstrap.service display-manager.service

[Service]
Type=oneshot
ExecStart=/usr/local/lib/nmos/logging_policy.py
StandardOutput=journal+console
StandardError=journal+console
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/etc/systemd/journald.conf.d @NMOS_RUNTIME_DIR@
CapabilityBoundingSet=
RestrictAddressFamilies=AF_UNIX

[Install]
WantedBy=multi-user.target
