# config/system-overlay/usr/lib/systemd/system/nmos-update-health.service
[Unit]
Description=NM-OS update boot health monitor
After=nmos-update-engine.service
Wants=nmos-update-engine.service

[Service]
Type=simple
ExecStart=/usr/local/lib/nmos/update_boot_health.py
StandardOutput=journal+console
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=@NMOS_RUNTIME_DIR@ @NMOS_STATE_DIR@
CapabilityBoundingSet=
RestrictAddressFamilies=AF_UNIX
Restart=no

[Install]
WantedBy=multi-user.target
