# config/system-overlay/usr/lib/systemd/system/nmos-persistent-storage.service
[Unit]
Description=NM-OS persistent storage backend
After=local-fs.target dbus.service
Wants=dbus.service

[Service]
Type=simple
ExecStart=/usr/local/bin/nmos-persistent-storage
Restart=on-failure
RestartSec=3
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=@NMOS_RUNTIME_DIR@ @NMOS_STATE_DIR@
CapabilityBoundingSet=CAP_SYS_ADMIN CAP_DAC_OVERRIDE CAP_DAC_READ_SEARCH CAP_FOWNER CAP_CHOWN CAP_MKNOD CAP_SYS_RAWIO

[Install]
WantedBy=multi-user.target
